# systemd/rektrace-dev.service
[Unit]
Description=Rektrace Development Environment
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/Rektrace
ExecStart=%h/Rektrace/scripts/dev.sh start
ExecStop=%h/Rektrace/scripts/dev.sh stop
TimeoutStartSec=300
TimeoutStopSec=60

# Environment
Environment=DOCKER_BUILDKIT=1
Environment=COMPOSE_DOCKER_CLI_BUILD=1

# Restart policy
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target